# projects/Rockchip/devices/RK356X/packages/RK356X-volkeys/system.d/RK356X-volkeys.service
[Unit]
Description=Rockchip RK356X series series volume keys sense service
Before=retroarch.target

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/spkeys-service gpio-keys-vol
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
