# projects/Rockchip/devices/RK3326/packages/odroidgo2-utils/system.d/odroidgoa-utils.service
[Unit]
Description=Odroid Go Advance headphone sense
Before=retroarch.target

[Service]
Environment=HOME=/storage
ExecStart=/usr/bin/headphone-sense 'rk817_int Headphones'
Restart=on-failure

[Install]
WantedBy=multi-user.target
